# profile.d/rednet.sh

for script in /etc/bash_completion.d/*
do
    . $script
done

# there was a script to set up apiserver env here
